# 07/challenges/06/apache2.conf
[program:apache2]
command=/usr/sbin/apache2ctl -D "FOREGROUND" -k start
autostart=true
autorestart=true
startretries=3
stopwaitsecs=10